# .cache/crontab/crontab.bak
0 */6 * * * /data/data/com.termux/files/home/run-auto-resolve.sh
